# server/files/etc/cron.d/backup_hbase_to_s3
#   Full snapshot every 2 hours
20 */2 * * *  root  test -f /etc/ec2onrails/ebs_hbase_info.yml && /usr/local/ec2onrails/bin/backup_hbase.rb
